# trunk/vendors/Nintendo/DLDI/inittab
tty1::linux:/usr/bin/agetty -n -l /bin/autologin 38400 tty1

