# branches/stable/dslinux_2_6_9/vendors/Nintendo/DSGBA/inittab

tty1::linux:/bin/agetty 38400 tty1
tty2::linux:/bin/agetty 38400 tty2
tty3::linux:/bin/agetty 38400 tty3
tty4::linux:/bin/agetty 38400 tty4

