# trunk/vendors/Nintendo/RAM-DEVEL/inittab

tty1::linux:/usr/bin/agetty 38400 tty1

