# vendors/Nintendo/RAM/inittab
tty1::linux:/usr/bin/agetty -n -l /bin/autologin 38400 tty1
tty2::linux:/usr/bin/agetty -n -l /bin/autologin 38400 tty2
tty3::linux:/usr/bin/agetty -n -l /bin/autologin 38400 tty3
